# package/base-files/files/etc/hotplug.d/ntp/00-syncrtc
#!/bin/sh
logger -t ntp-plug "ntp plug shell,$ACTION"
/sbin/hwclock -w
