# package/utils/net4g/files/media_resart
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

pidof mediasvr && killall -INT mediasvr
killall -9 mediasvr
/bin/mediasvr &
