# package/utils/net4g/files/net4g.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=80

start() {
	pidof net4g || /bin/net4g -d
}

stop() {
	killall -INT net4g
}
